# crontab
10 01 * * SUN /bin/sh /cronscript.sh >> cronoutput.txt
